# deploy/supervisor_redis.conf

[program:redis]
user=root
command=redis-server /etc/redis/redis.conf
autorestart=false
autostart=true

